# supervisord.ini
[supervisord]
nodaemon=true

[program:django]
command=bash src/start.sh
stderr_logfile=/var/log/django.err.log
stderr_logfile_maxbytes=0
stdout_logfile=/var/log/django.out.log
stdout_logfile_maxbytes=0
stopsignal=QUIT
environment=PYTHONUNBUFFERED=1

[program:celery]
command=celery -A megano worker --loglevel=info
stderr_logfile=/var/log/celery.err.log
stderr_logfile_maxbytes=0
stdout_logfile=/var/log/celery.out.log
stdout_logfile_maxbytes=0
stopsignal=QUIT
environment=PYTHONUNBUFFERED=1
