# supervisord.ini
[supervisord]
nodaemon=true

[program:django]
command=bash src/start.sh
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stopsignal=QUIT
environment=PYTHONUNBUFFERED=1

[program:celery]
command=celery -A megano worker --loglevel=info
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stopsignal=QUIT
environment=PYTHONUNBUFFERED=1
